# cron/crontab
TZ=Asia/Tokyo

05 * * * * source /root/suna/envfile ; cd / ; /repo/bin/init >> /var/log/cron.log 2>&1
* * * * * echo "aaaaaaaaaaaa" > /proc/1/fd/1 2>&1
