# cron/crontab
TZ=Asia/Tokyo

05 * * * * /cron/exec-init.sh > /proc/1/fd/1 2>&1
00 * * * * date > /proc/1/fd/1 2>&1
