# docker/cron/crontab
05 * * * * cd / ; /init >> /var/log/cron.log 2>&1