# docker/cron/crontab
05 * * * * source /root/suna/envfile ; cd / ; /init >> /var/log/cron.log 2>&1
* * * * * echo "aaaaaaaaaaaa" >> /var/log/cron.log 2>&1
